# supervisor/bsky_feed.conf
[program:bsky_feed]
directory=/root/bsky-feed
command=yarn start
autostart=true
autorestart=true

